# test.service
[Unit]
Description=A bit about your service

[Service] 
Type=simple
ExecStart=/usr/bin/wall "hello, whirled!"

[Install]
WantedBy=multi-user.target
